# ansible/07-kube-controller-manager/templates/kube-controller-manager.service.j2
[Unit]
Description=Kubernetes Controller Manager
Documentation=https://kubernetes.io/docs/reference/command-line-tools-reference/kube-controller-manager/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kubernetes
Group=kubernetes

WorkingDirectory=/var/lib/kubernetes
EnvironmentFile=-/etc/kubernetes/kube-controller-manager.env
ExecStart=/usr/local/bin/kube-controller-manager $KUBE_CONTROLLER_MANAGER_OPTS

Restart=on-failure
RestartSec=10

# Limites de recursos
LimitNOFILE=65536
LimitNPROC=65536

# Configurações de segurança
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/kubernetes

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kube-controller-manager

[Install]
WantedBy=multi-user.target