# ansible/kube_apiserver/templates/kube-apiserver.service.j2
[Unit]
Description=Kubernetes API Server
Documentation=https://kubernetes.io/docs/reference/command-line-tools-reference/kube-apiserver/
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=kubernetes
Group=kubernetes

ExecStart=/usr/local/bin/kube-apiserver \
  --advertise-address={{ hostvars[groups['loadbalancers'][0]]['ansible_host'] }} \
  --allow-privileged=true \
  --audit-log-maxage=30 \
  --audit-log-maxbackup=3 \
  --audit-log-maxsize=100 \
  --audit-log-path=/var/log/kubernetes/audit.log \
  --authorization-mode=Node,RBAC \
  --bind-address=0.0.0.0 \
  --client-ca-file=/etc/kubernetes/pki/{{ certs.ca_k8s.cert | basename }} \
  --enable-admission-plugins=NamespaceLifecycle,NodeRestriction,LimitRanger,ServiceAccount,DefaultStorageClass,ResourceQuota \
  --etcd-cafile=/etc/kubernetes/pki/{{ certs.ca_etcd.cert | basename }} \
  --etcd-certfile=/etc/kubernetes/pki/apiserver/{{ certs.etcd_apiserver_cliente.cert | basename }} \
  --etcd-keyfile=/etc/kubernetes/pki/apiserver/{{ certs.etcd_apiserver_cliente.chave | basename }} \
  --etcd-servers=https://{{ hostvars[groups['loadbalancers'][0]]['ansible_host'] }}:{{ porta_etcd_https }} \
  --etcd-healthcheck-timeout=5s \
  --etcd-compaction-interval=5m \
  --etcd-count-metric-poll-period=1m \
  --event-ttl=1h \
  --kubelet-certificate-authority=/etc/kubernetes/pki/{{ certs.ca_k8s.cert | basename }} \
  --kubelet-client-certificate=/etc/kubernetes/pki/apiserver/{{ certs.kube_apiserver.cert | format(inventory_hostname) | basename }} \
  --kubelet-client-key=/etc/kubernetes/pki/apiserver/{{ certs.kube_apiserver.chave | format(inventory_hostname) | basename }} \
  --max-requests-inflight={{ max_requests | default(1500) }} \
  --requestheader-client-ca-file=/etc/kubernetes/pki/{{ certs.ca_k8s.cert | basename }} \
  --requestheader-allowed-names=front-proxy-client \
  --requestheader-extra-headers-prefix=X-Remote-Extra- \
  --requestheader-group-headers=X-Remote-Group \
  --requestheader-username-headers=X-Remote-User \
  --runtime-config=api/all=true \
  --service-account-key-file=/etc/kubernetes/pki/apiserver/{{ certs.service_account.cert | basename }} \
  --service-account-signing-key-file=/etc/kubernetes/pki/apiserver/{{ certs.service_account.chave | basename }} \
  --service-account-issuer=https://kubernetes.default.svc.cluster.local \
  --service-cluster-ip-range={{ rede_cidr_services }} \
  --service-node-port-range=30000-32767 \
  --tls-cert-file=/etc/kubernetes/pki/apiserver/{{ certs.kube_apiserver.cert | format(inventory_hostname) | basename }} \
  --tls-private-key-file=/etc/kubernetes/pki/apiserver/{{ certs.kube_apiserver.chave | format(inventory_hostname) | basename }} \
  --api-audiences=https://kubernetes.default.svc.cluster.local \
  --enable-bootstrap-token-auth=true \
  --proxy-client-cert-file=/etc/kubernetes/pki/apiserver/{{ certs.kube_apiserver.cert | format(inventory_hostname) | basename }} \
  --proxy-client-key-file=/etc/kubernetes/pki/apiserver/{{ certs.kube_apiserver.chave | format(inventory_hostname) | basename }} \
  --anonymous-auth=false \
  --secure-port={{ porta_kubeapiserver }} \
  --v=2

Restart=on-failure
RestartSec=5
TimeoutStartSec=0
TimeoutStopSec={{ timeout_shutdown | default(70) }}

# Limites de recursos
LimitNOFILE=65536
LimitNPROC=65536

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kube-apiserver

# Configurações de segurança
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/kubernetes /var/log/kubernetes

[Install]
WantedBy=multi-user.target